# app-misc/rogentoslive-tools/files/2/installer-text.service
[Unit]
Description=Sabayon text installer setup
After=rogentoslive.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/libexec/installer-text.sh

[Install]
WantedBy=multi-user.target
